# files/rc.local
#!/bin/bash

/etc/init.d/procps restart
screen -AmdS ohp-dropbear ohpserver -port 3128 -proxy 127.0.0.1:8080 -tunnel 127.0.0.1:110
screen -AmdS ohp-openvpn ohpserver -port 8000 -proxy 127.0.0.1:8080 -tunnel 127.0.0.1:992
screen -AmdS badvpn badvpn-udpgw --listen-addr 127.0.0.1:7300
iptables-restore < /metavpn/iptables.rules

exit 0